# aci-assets/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:postfix]
command=/usr/local/bin/postfix_start.sh
